# services/dhcpcd/run
#!/bin/sh

exec /usr/sbin/dhcpcd -B eth0
