# services/dbus/run
#!/bin/sh
sleep 5
exec dbus-daemon --system --nofork --nopidfile
