# services/dbus/run
#!/bin/sh
sleep 3
exec dbus-daemon --system --nofork --nopidfile
