# services/eudev/finish
#!/bin/sh

/usr/sbin/udevadm trigger
/usr/sbin/udevadm settle