# services/eudev/finish
#!/bin/sh

/sbin/udevadm trigger
/sbin/udevadm settle