# services/udev-input-trigger/run
#!/bin/sh
sleep 5
for dev in /sys/class/input/event*; do
    devpath=$(readlink -f "$dev/device")
    udevadm test "$devpath" > /dev/null 2>&1
done
touch /etc/sv/udev-input-trigger/down
