# services/console/run
#!/bin/sh

exec setsid sh -c 'exec sh </dev/ttyS0 >/dev/ttyS0 2>&1'