# services/console/run
#!/bin/sh

exec setsid sh -c 'exec /sbin/getty -a root 38400 ttyS0 </dev/ttyS0 >/dev/ttyS0 2>&1'